# githubback.service
[Unit]
Description=MyScript

[Service]
ExecStart=/home/liam/backup-github/githuback
